# src/scripts/systemd/nfs-ganesha-config.service
[Unit]
Description=Process NFS-Ganesha configuration
Before=nfs-ganesha.service

[Service]
ExecStart=/usr/libexec/ganesha/nfs-ganesha-config.sh
